# lupinus.ini
[program:lupinus]
command=/var/www/lupinus/start.sh
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/jobs/lupinus-supervisord.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true